# stow/dot-local/service/syncthing/run
#!/usr/bin/env sh

exec syncthing --no-browser
