# stow/home/dot-local/service/picom/run
#!/usr/bin/env dash

exec picom --backend xrender

