# stow/home/dot-local/service/picom/run
#!/bin/sh

exec picom --backend xrender

