# stow/home/dot-local/service/shadowsocks-no/run
#!/bin/sh

exec sslocal -c $HOME/.local/service/shadowsocks-no/config.json
