# stow/home/dot-local/service/shadowsocks-no/run
#!/bin/sh

exec ss-local -c $HOME/.local/service/shadowsocks-no/config.json
