# stow/home/dot-local/service/sxhkd/run
#!/bin/sh

common="$HOME/.config/sxhkd/sxhkdrc"

if [ "$DESKTOP_SESSION" = "none+mind-shift" ]; then
    extra="$HOME/.config/sxhkd/sxhkdrc-ms"
fi

exec sxhkd -c $common $extra
