# stow/home/dot-local/service/mpd/run
#!/usr/bin/env dash

exec mpd --no-daemon ~/.config/mpd/mpd.conf

