# stow/home/dot-local/service/mpd/run
#!/bin/sh

exec mpd --no-daemon ~/.config/mpd/mpd.conf

