# dot-local/service/shadowsocks-nl/run
#!/usr/bin/env dash

exec sslocal -c $HOME/.config/shadowsocks/config-nl.json
