# stow/home/dot-local/service/shadowsocks-nl/run
#!/usr/bin/env dash

exec sslocal -c $HOME/.local/service/shadowsocks-nl/config.json
