# stow/home/dot-local/service/shadowsocks-nl/run
#!/bin/sh

exec ss-local -c $HOME/.local/service/shadowsocks-nl/config.json
