# dot-local/service/win-layout/run
#!/usr/bin/env dash

exec $HOME/.local/bin/win-layout
