# dot-local/service/xbindkeys/run
#!/usr/bin/env dash

exec xbindkeys -n -f $HOME/.local/service/xbindkeys/xbindkeysrc
